# supervisord/filebeat.conf
[program:filebeat]
user=root
command=/usr/share/filebeat/bin/filebeat
	-c /etc/filebeat/filebeat.yml
	-path.home /usr/share/filebeat
	-path.config /etc/filebeat
	-path.data /var/lib/filebeat
	-path.logs /var/log/filebeat
autostart=true
autorestart=false
stdout_logfile=/var/log/supervisor/filebeat.out.log
stderr_logfile=/var/log/supervisor/filebeat.err.log
