# supervisord/rsyslog.conf
[program:rsyslog]
user=root
command=/usr/sbin/rsyslogd -n
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/rsyslog.out.log
stderr_logfile=/var/log/supervisor/rsyslog.err.log
